# deployment/supervisord.dev-hot.conf
# Development supervisord config with hot-reload support
# Used by docker-compose.dev-hot.yml for fast iteration.
#
# Key differences from production supervisord.conf:
# - uvicorn runs with --reload to watch for Python file changes (~1-2s restart)
# - Next.js runs in dev mode (next dev --turbopack) for instant HMR
#
# This is the lean variant: nginx + fastapi + nextjs only. Downstream
# deployments that need further programs (e.g. a VPN sidecar or telemetry
# forwarder) should add them in their own copy; keep those out of here so
# supervisord does not reference env vars the community compose never sets.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3

[program:fastapi]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --log-level info --reload --reload-dir /app/app
directory=/app
priority=6
autostart=true
autorestart=true
# uvicorn --reload spawns a reloader + worker child; stop/kill the whole group
# so a supervisor restart doesn't leave straggler processes in the dev loop.
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
startretries=3
stopwaitsecs=10

[program:nextjs]
command=/bin/sh -c 'NODE_ENV=development PORT=3000 exec npm run dev'
directory=/app/imi-frontend
priority=7
autostart=true
autorestart=true
# next dev forks child processes; group stop/kill avoids orphaned node procs.
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=3
stopwaitsecs=10
